# etc/downlink.service
[Unit]
Description=DOWNLINK Feed Aggregator Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/downlink
ExecStart=/opt/downlink/server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=downlink

[Install]
WantedBy=multi-user.target
